# hyprland-autoname-workspaces.service
[Unit]
Description=Autoname Hyprland workspaces
PartOf=graphical-session.target
After=graphical-session.target

[Service]
ExecStart=/usr/bin/hyprland-autoname-workspaces
Restart=always
RestartSec=10s

[Install]
WantedBy=graphical-session.target
